# tests/pkcli/build_data/1-fconf.out/srv/host/v9.radia.run/etc/systemd/system/jupyterhub.service
[Unit]
After=docker.service 
Description=jupyterhub
Requires=docker.service

[Service]
ExecStart=/srv/jupyterhub/start
ExecStartPre=/srv/jupyterhub/remove
ExecStop=-/srv/jupyterhub/stop
Restart=on-failure
RestartSec=10
SyslogIdentifier=jupyterhub

[Install]
WantedBy=multi-user.target
