# tests/pkcli/build_data/1-fconf.out/srv/host/v9.radia.run/etc/systemd/system/sirepo_test_http.timer
[Unit]
Description=sirepo_test_http

[Timer]
OnCalendar=*-*-* 15:00:0

[Install]
WantedBy=timers.target
