# tests/pkcli/build_data/1.out/srv/host/v9.radia.run/etc/systemd/system/btest.service
[Unit]
Description=btest

[Service]
Environment=TZ=:/etc/localtime
ExecStart=/srv/btest/start
Group=vagrant
SyslogIdentifier=btest
Type=oneshot
User=vagrant

[Install]
WantedBy=multi-user.target
