# tests/pkcli/build_data/1.out/srv/host/v9.radia.run/etc/systemd/system/rs_mariadb.service
[Unit]
After=docker.service 
Description=rs_mariadb
Requires=docker.service

[Service]
ExecStart=/srv/rs_mariadb/start
ExecStartPre=/srv/rs_mariadb/remove
ExecStop=-/srv/rs_mariadb/stop
SyslogIdentifier=rs_mariadb
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
