# tests/pkcli/build_data/1.out/srv/host/v9.radia.run/etc/systemd/system/sirepo_test_http.service
[Unit]
After=docker.service 
Description=sirepo_test_http
Requires=docker.service

[Service]
ExecStart=/srv/sirepo_test_http/start
ExecStartPre=/srv/sirepo_test_http/remove
ExecStop=-/srv/sirepo_test_http/stop
SyslogIdentifier=sirepo_test_http
Type=oneshot

[Install]
WantedBy=multi-user.target
