# tests/pkcli/build_data/1-jinja.out/srv/host/v5.radia.run/etc/systemd/system/logrotate.timer
[Unit]
Description=logrotate

[Timer]
OnCalendar=*-*-* 13:30:0

[Install]
WantedBy=timers.target
