# tests/pkcli/build_data/1-jinja.out/srv/host/v5.radia.run/etc/systemd/system/mpi_worker.service
[Unit]
After=docker.service 
Description=mpi_worker
Requires=docker.service

[Service]
ExecStart=/srv/mpi_worker/start
ExecStartPre=/srv/mpi_worker/remove
ExecStop=-/srv/mpi_worker/stop
Restart=on-failure
RestartSec=10
SyslogIdentifier=mpi_worker

[Install]
WantedBy=multi-user.target
