# tests/pkcli/build_data/jinja1.out/srv/host/v4.radia.run/etc/systemd/system/db_bkp.timer
[Unit]
Description=db_bkp

[Timer]
OnCalendar=*-*-* 12:30:0

[Install]
WantedBy=timers.target
