# tests/pkcli/build1_data/1.out/srv/host/v3.radia.run/etc/systemd/system/db_bkp.timer
[Unit]
Description=db_bkp

[Timer]
OnCalendar=*-*-* 13:30:0

[Install]
WantedBy=timers.target
