# tests/pkcli/build1_data/1.out/srv/host/v3.radia.run/etc/systemd/system/logrotate.service
[Unit]
Description=logrotate

[Service]
Environment=TZ=:/etc/localtime
ExecStart=/srv/logrotate/start
Group=root
SyslogIdentifier=logrotate
Type=oneshot
User=root
