# tests/pkcli/build_data/1-fconf.out/srv/host/v4.radia.run/etc/systemd/system/sirepo_jupyterhub.service
[Unit]
After=docker.service 
Description=sirepo_jupyterhub
Requires=docker.service

[Service]
ExecStart=/srv/sirepo_jupyterhub/start
ExecStartPre=/srv/sirepo_jupyterhub/remove
ExecStop=-/srv/sirepo_jupyterhub/stop
SyslogIdentifier=sirepo_jupyterhub
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
