# tests/pkcli/build_data/1-jinja.out/srv/host/v9.radia.run/etc/systemd/system/github_bkp.timer
[Unit]
Description=github_bkp

[Timer]
OnCalendar=*-*-* 9:0:0

[Install]
WantedBy=timers.target
