# tests/pkcli/build1_data/1.out/srv/host/v9.radia.run/etc/systemd/system/bkp.timer
[Unit]
Description=bkp

[Timer]
OnCalendar=*-*-1 23:59:0

[Install]
WantedBy=timers.target
