# tests/pkcli/build1_data/1.out/srv/host/v9.radia.run/etc/systemd/system/bkp.service
[Unit]
Description=bkp

[Service]
Environment=TZ=:/etc/localtime
ExecStart=/srv/bkp/start
Group=root
SyslogIdentifier=bkp
Type=oneshot
User=root

[Install]
WantedBy=multi-user.target
