# tests/pkcli/build1_data/1.out/srv/host/v5.radia.run/etc/systemd/system/db_bkp.timer
[Unit]
Description=db_bkp

[Timer]
OnCalendar=*-*-* 6:30:0 America/Phoenix

[Install]
WantedBy=timers.target
