# tests/pkcli/build1_data/1.out/srv/host/v5.radia.run/etc/systemd/system/logrotate.timer
[Unit]
Description=logrotate

[Timer]
OnCalendar=*-*-* 7:30:0 America/Phoenix

[Install]
WantedBy=timers.target
